# autologin.service
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin displayuser --noclear %I 38400 linux
